# packages/net-fs/cachefilesd/files/systemd/cachefilesd.service
[Unit]
Description=FSCache cleanup daemon
Wants=syslog.target

[Service]
ExecStart=/sbin/cachefilesd -n

[Install]
WantedBy=remote-fs.target
