# packages/sys-apps/mlocate/files/systemd/mlocate.service
[Unit]
Description=mlocate updatedb

[Service]
Type=oneshot
ExecStart=/usr/bin/updatedb
