# packages/sys-apps/mlocate/files/systemd/mlocate.timer
[Unit]
Description=mlocate updatedb timer

[Timer]
OnUnitActiveSec=1d

[Install]
WantedBy=multi-user.target
