# packages/www-servers/apache/files/systemd/apache.service
[Unit]
Description=apache2
After=network.target

[Service]
Type=forking
ExecStartPre=/usr/bin/apachectl -t
ExecStart=/usr/bin/apachectl -k start
ExecStop=/usr/bin/apachectl -k stop
ExecReload=/usr/bin/apachectl -k graceful

[Install]
WantedBy=http-daemon.target
