# packages/net-fs/moosefs/files/systemd/mfsmaster.service
[Unit]
Description=MooseFS master node
After=syslog.target network.target

# This daemon forks out only after it's successfully
#  initialized, so there's an advantage in running it like this
[Service]
Type=forking
TimeoutSec=600
ExecStartPre=/bin/su -s /usr/sbin/mfsprepare mfs
ExecStart=/usr/sbin/mfsmaster start
ExecStop=/usr/sbin/mfsmaster stop

[Install]
WantedBy=multi-user.target
