# packages/net-fs/moosefs/files/systemd/mfsmetalogger.service
[Unit]
Description=MooseFS metalogger
After=syslog.target network.target mfsmaster.service

# This daemon forks out only after it's successfully
#  initialized, so there's an advantage in running it like this
[Service]
Type=forking
ExecStart=/usr/sbin/mfsmetalogger start
ExecStop=/usr/sbin/mfsmetalogger stop

[Install]
WantedBy=multi-user.target
