# packages/net-proxy/polipo/files/systemd/polipo.service
[Unit]
Description=polipo
After=network.target

[Service]
ExecStart=/usr/bin/polipo -c /etc/polipo.conf

[Install]
WantedBy=multi-user.target
