# packages/net-proxy/ziproxy/files/systemd/ziproxy.service
[Unit]
Description=ziproxy
After=network.target

[Service]
Type=forking
PIDFile=/run/ziproxy.pid
ExecStart=/usr/bin/ziproxy -d --pid-file=/run/ziproxy.pid
ExecStop=/usr/bin/ziproxy -k --pid-file=/run/ziproxy.pid

[Install]
WantedBy=multi-user.target
