# packages/net-proxy/ziproxy/files/systemd/ziproxy.service
[Unit]
Description=ziproxy
Requires=network.target
After=network.target

[Service]
Type=forking
PIDFile=/var/run/ziproxy.pid
ExecStart=/usr/bin/ziproxy -d --pid-file=/var/run/ziproxy.pid
ExecStop=/usr/bin/ziproxy -k --pid-file=/var/run/ziproxy.pid

[Install]
WantedBy=multi-user.target
