# packages/dev-libs/openct/files/systemd/openct-handler@.service
[Unit]
Requires=openct.service

[Service]
Type=forking
GuessMainPID=false
ExecStart=/bin/sh -c "openct-control attach %I"
