# packages/monitor/net-snmp/files/systemd/snmpd.service
[Unit]
Description=snmpd

[Service]
ExecStart=/usr/sbin/snmpd -f -u snmpd -g snmpd
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
