# packages/dev-db/mongodb/files/systemd/mongod.service
[Unit]
Description=mongod

[Service]
EnvironmentFile=/etc/conf.d/mongodb.conf
User=$MONGOD_USER
Group=$MONGOD_GROUP
ExecStart=/usr/bin/mongod $MONGOD_OPTIONS

[Install]
WantedBy=multi-user.target
