# packages/net-fs/glusterfs/files/systemd/glusterd.service
[Unit]
Description=glusterd

[Service]
User=glusterd
Group=glusterd
Capabilities=cap_sys_admin=i
ExecStart=/usr/sbin/glusterd --no-daemon

[Install]
WantedBy=multi-user.target
