# packages/net-fs/glusterfs/files/systemd/glusterd.service
[Unit]
Description=glusterd

[Service]
ExecStart=/usr/sbin/glusterd --no-daemon

[Install]
WantedBy=multi-user.target
