# packages/net-dns/djbdns/files/systemd/dnscache.service
[Unit]
Description=dnscache
Requires=network.target
After=network.target

[Service]
EnvironmentFile=/var/lib/djbdns/cache/env
ExecStart=/bin/sh -c 'exec /usr/bin/dnscache </var/lib/djbdns/cache/seed'

[Install]
WantedBy=multi-user.target
