# packages/net-dns/djbdns/files/systemd/dnscache.service
[Unit]
Description=dnscache
Requires=network.target
After=network.target

# logger here is a hack, but systemd doesn't seem to be able to log from non-kernel facility
# TODO: use proper systemd-logger capabilities (as per systemd.exec(5)), when it'll be fixed
[Service]
EnvironmentFile=/etc/djbdns/cache/env
ExecStart=/bin/sh -c '/usr/bin/dnscache </etc/djbdns/cache/seed 2>&1 | /usr/bin/logger -t dnscache -p daemon.info'

[Install]
WantedBy=multi-user.target
