# packages/net-dns/djbdns/files/systemd/tinydns.service
[Unit]
Description=tinydns
Requires=network.target
After=network.target

[Service]
EnvironmentFile=/etc/djbdns/resolver/env
StandardOutput=null
StandardError=inherit
SyslogFacility=daemon
SyslogLevel=info
SyslogIdentifier=tinydns
ExecStart=/usr/bin/tinydns

[Install]
WantedBy=multi-user.target
