# packages/net-dns/djbdns/files/systemd/tinydns.service
[Unit]
Description=tinydns
Requires=network.target
After=network.target

[Service]
EnvironmentFile=/etc/djbdns/resolver/env
StandardOutput=syslog
StandardError=syslog
SyslogFacility=daemon
SyslogLevel=info
ExecStart=/usr/bin/tinydns

[Install]
WantedBy=multi-user.target
