# packages/dev-lang/erlang/files/systemd/epmd.service
[Unit]
Description=Erlang Portmapper Daemon
Before=network.target

[Service]
ExecStart=/usr/bin/epmd

[Install]
WantedBy=multi-user.target
