# packages/monitor/collectd/files/systemd/collectd.service
[Service]
ExecStart=/usr/sbin/collectd -f

[Install]
WantedBy=multi-user.target
