# packages/monitor/graphite-carbon/files/systemd/graphite_carbon_relay@.service
[Unit]
Description=graphite carbon relay
After=network.target

[Service]
Type=forking
User=graphite
GuessMainPID=false
PIDFile=/opt/graphite/storage/carbon-relay-%i.pid
ExecStartPre=-/bin/rm -f /opt/graphite/storage/carbon-relay-%i.pid
ExecStart=/opt/graphite/bin/carbon-relay.py --pidfile=/opt/graphite/storage/carbon-relay-%i.pid --instance=%i start
ExecStop=/opt/graphite/bin/carbon-relay.py --pidfile=/opt/graphite/storage/carbon-relay-%i.pid --instance=%i stop

[Install]
WantedBy=multi-user.target
