# packages/monitor/shinken/files/systemd/shinken@.service
[Unit]
Description=Shinken %i

[Service]
WorkingDirectory=/var/lib/shinken
SyslogIdentifier=shinken-%i
ExecStart=/usr/bin/env shinken-%i -c /etc/shinken/%id.ini

[Install]
WantedBy=multi-user.target
