# packages/monitor/ganglia/files/systemd/gmond.service
[Unit]
Description=gmond
After=network.target

[Service]
ExecStart=/usr/sbin/gmond -f

[Install]
WantedBy=multi-user.target
