# packages/net-im/ejabberd/files/systemd/ejabberd.service
[Unit]
Description=ejabberd
After=network.target epmd.service

[Service]
Type=forking
User=jabber
Group=jabber
GuessMainPID=false
ExecStart=/usr/sbin/ejabberdctl start
ExecStop=/usr/sbin/ejabberdctl stop

[Install]
WantedBy=multi-user.target
