# packages/net-im/ejabberd/files/systemd/ejabberd.service
[Unit]
Description=ejabberd
Requires=network.target
After=network.target

[Service]
Type=forking
User=jabber
Group=jabber
ExecStart=/usr/sbin/ejabberdctl start
ExecStop=/usr/sbin/ejabberdctl stop

[Install]
WantedBy=multi-user.target
