# packages/net/sslh/files/systemd/sslh.service
[Unit]
Description=sslh
After=syslog.target

[Service]
ExecStart=/usr/sbin/sslh -u nobody -p 0.0.0.0:443 -s 127.0.0.1:22 -l 127.0.0.1:8102

[Install]
WantedBy=multi-user.target
