# packages/net-fs/cfs/files/systemd/cfs.service
[Unit]
Description=cfs

[Service]
ExecStart=/usr/sbin/cfs -c /etc/cfs.yaml

[Install]
WantedBy=multi-user.target
