# packages/net-fs/cfs/files/systemd/cfsd.service
[Unit]
Description=cfsd

[Service]
ExecStart=/usr/sbin/cfsd -c /etc/cfsd.yaml

[Install]
WantedBy=multi-user.target
